# systemd/robomesha.service
[Unit]
Description=ROBOMESHA stack
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=__ROOT_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__ROOT_DIR__/scripts/start_stack.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
